# sdcard/rootfs/docker-multinode/etc/profile.d/docker-bootstrap.sh
system-docker(){
    docker -H unix:///var/run/docker-bootstrap.sock "$@"
}
docker-bootstrap(){
    docker -H unix:///var/run/docker-bootstrap.sock "$@"
}
docker-rm-stopped(){ 
    docker rm $(docker ps --filter status=exited -q) 
}
